# udev/99-backive.rules
KERNEL=="sd*", RUN+="/path/to/backive_udev"
